# etc/systemd/system/docker.service.d/TimeoutStopSec.conf
# /etc/systemd/system/docker.service.d/TimeoutStopSec.conf

[Service]
# Set the appropriate value according to your guest OS needs
TimeoutStopSec=5h
